# cloudflared.service
[Unit]
Description=Cloudflare Tunnel for Pool Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
ExecStart=/usr/local/bin/cloudflared tunnel run
Restart=on-failure
RestartSec=5s
StandardOutput=append:/var/log/cloudflared.log
StandardError=append:/var/log/cloudflared-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
